# share/cron.d/ardj
# /etc/cron.d/ardj: crontab fragment for ardj
@daily ardj ardj mark-hitlist mark-recent mark-orphans mark-long update-schedule purge
