# build/rootfs/etc/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock

[supervisord]
logfile = /var/log/supervisord.log
pidfile = /run/supervisord.pid
childlogdir = /var/log
nodaemon = true
user = root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:icecast]
command = /usr/bin/icecast -c /app/data/icecast2.xml
directory = /app
autostart = true
autorestart = true
startsecs = 1
numprocs = 1
startretries = 1000
restartpause = 5
user = root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
priority = 1

[program:ezstream]
command = /usr/bin/ezstream -v -c /app/data/ezstream.xml
directory = /app
autostart = true
autorestart = true
startsecs = 1
numprocs = 1
#process_name = %(process_num)03d
startretries = 1000
restartpause = 5
user = radio
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
priority = 2

[program:server]
command = python3 -m ardj serve
directory = /app
autostart = true
autorestart = true
startsecs = 1
numprocs = 1
#process_name = %(process_num)03d
startretries = 1000
restartpause = 5
user = radio
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:scrobbler]
command = python3 -m ardj scrobbler
directory = /app
autostart = true
autorestart = true
startsecs = 1
numprocs = 1
#process_name = %(process_num)03d
startretries = 1000
restartpause = 5
user = radio
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
